# usr/lib/systemd/system/akbl.service
[Unit]
Description=AKBL Daemon

[Service]
ExecStart=/bin/bash -c "/usr/bin/akbl --start-daemon"

[Install]
WantedBy=basic.target
