# usr/lib/systemd/system/alienware-kbl.service
[Unit]
Description=alienware-kbl daemon

[Service]
ExecStart=/bin/bash -c "/usr/bin/alienware-kbl --start-daemon"

[Install]
WantedBy=basic.target
